# rootfs/etc/s6/bind/run
#!/bin/sh

exec named -p 53 -u named -c /etc/bind/named.conf -g
